# supervisord.conf
[supervisord]
nodaemon=true

[program:parse]
command=npm start
